# files/rc.local
#!/bin/sh -e
# This file is managed by Puppet.

sleep 10 && bash /opt/newts/bin/newts -D -c /opt/newts/etc/config.yaml 2>&1 > /tmp/newts_start.log &

exit 0

